# crontab
*/15 * * * * root urlwatch -v >> /var/log/urlwatch.log 2>&1
* * * * * root urlwatch -v >> /var/log/urlwatch2.log 2>&1
0 2  * * *  root apt-get update && apt-get dist-upgrade -y
