# crontab
*/15 * * * * /usr/bin/urlwatch -v >> /var/log/urlwatch.log 2>&1
0 2  * * *  /usr/bin/apt-get update && /usr/bin/apt-get dist-upgrade -y
0 12 * * * echo "Start" > /var/log/urlwatch.log
